# etc/profile.d/enforce-platform.sh
export PLASMA_PLATFORM=mediacenter
export PLASMA_USE_QT_SCALING=true
export QT_SCREEN_SCALE_FACTORS=1.5
export BIGSCREEN_HARDWARE_PLATFORM=RPI4
