# etc/profile.d/enforce-platform.sh
export PLASMA_PLATFORM=mediacenter
export BIGSCREEN_HARDWARE_PLATFORM=RPI4
